# provisioners/docker/images/traefik/supervisord.tmpl.conf
[unix_http_server]
file=/tmp/supervisor.sock                       ; path to your socket file

[supervisord]
logfile=/var/log/supervisord/supervisord.log    ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=debug                                  ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid                ; pidfile location
nodaemon=true                                   ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=root                                       ; default user
childlogdir=/var/log/supervisord/               ; where child log files will live

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock           ; use a unix:// URL for a unix socket

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:traefik]
command=/usr/local/bin/traefik
autorestart=true
redirect_stdout=true
redirect_stderr=true

{{if .Env.SUPERVISOR_ENABLE_BACKUP}}
[program:cron]
command=/usr/sbin/crond -f
autorestart=true
redirect_stdout=true
redirect_stderr=true

[program:acme-backup]
command=/usr/local/bin/backup.sh
autostart=false
autorestart=false
redirect_stdout=true
redirect_stderr=true
startsecs=0
{{end}}