# docker/images/traefik/backup.cron
# Run daily
0 0 * * *  /usr/local/bin/backup.sh
